# debian/debian_openport/postinst
systemctl enable openport
openport restart-sessions
